# ScheduledTaskPackage/TaskDefinition.xml
<?xml version="1.0" encoding="utf-8" ?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task"
      xmlns:xsi="http://www.w3.org/2001/XMLSchema-instance"
      xsi:schemaLocation="http://schemas.microsoft.com/windows/2004/02/mit/task ./TaskDefinition.xsd">
  <RegistrationInfo>
    <Author></Author>
    <Version></Version>
    <Description></Description>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Settings>
    <Enabled>true</Enabled>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <AllowHardTerminate>true</AllowHardTerminate>
  </Settings>
  <Actions>
    <Exec>
      <Command>$safeprojectname$.exe</Command>
      <Arguments></Arguments>
      <WorkingDirectory>#{Octopus.Action.Package.InstallationDirectoryPath}</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
